# build/rustdesk-server/files/rustdesk-template.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest"
                name="rustdesk">

    <service name="ooce/application/rustdesk"
             type="service"
             version="1">

        <dependency name="network"
                    grouping="require_all"
                    restart_on="error"
                    type="service">
            <service_fmri value="svc:/milestone/network:default" />
        </dependency>

        <dependency name="filesystem-local"
                    grouping="require_all"
                    restart_on="none"
                    type="service">
            <service_fmri value="svc:/system/filesystem/local:default" />
        </dependency>

        <instance name="hbbr"
                  enabled="false">

            <dependency name="hbbs"
                        grouping="require_all"
                        restart_on="restart"
                        type="service">
                <service_fmri value="svc:/ooce/application/rustdesk:hbbs" />
            </dependency>

            <exec_method type="method"
                         name="start"
                         exec="/opt/ooce/rustdesk-server/bin/hbbr &amp;"
                         timeout_seconds="60">
                <method_context security_flags="aslr">
                    <method_credential user="$(USER)"
                                       group="$(GROUP)"
                                       privileges="basic,!proc_info,!proc_session" />
                </method_context>
            </exec_method>

            <exec_method type="method"
                         name="stop"
                         exec=":kill"
                         timeout_seconds="60"></exec_method>

            <exec_method type="method"
                         name="refresh"
                         exec=":kill -HUP"
                         timeout_seconds="60"></exec_method>

            <property_group name="startd"
                            type="framework">
                <propval name="ignore_error"
                         type="astring"
                         value="core,signal" />
            </property_group>

        </instance>

        <instance name="hbbs"
                  enabled="false">

            <exec_method type="method"
                         name="start"
                         exec="/opt/ooce/rustdesk-server/bin/hbbs &amp;"
                         timeout_seconds="60">
                <method_context security_flags="aslr">
                    <method_credential user="$(USER)"
                                       group="$(GROUP)"
                                       privileges="basic,!proc_info,!proc_session" />
                </method_context>
            </exec_method>

            <exec_method type="method"
                         name="stop"
                         exec=":kill"
                         timeout_seconds="60"></exec_method>

            <exec_method type="method"
                         name="refresh"
                         exec=":kill -HUP"
                         timeout_seconds="60"></exec_method>

            <property_group name="startd"
                            type="framework">
                <propval name="ignore_error"
                         type="astring"
                         value="core,signal" />
            </property_group>

        </instance>

        <stability value="External" />

        <template>
            <common_name>
                <loctext xml:lang="C">RustDesk Server</loctext>
            </common_name>
        </template>

    </service>

</service_bundle>
