# build/nagios/files/nagios.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
	{{{ CDDL HEADER

	This file and its contents are supplied under the terms of the
	Common Development and Distribution License ("CDDL"), version 1.0.
	You may only use this file in accordance with the terms of version
	1.0 of the CDDL.
	A full copy of the text of the CDDL should have accompanied this
	source. A copy of the CDDL is also available via the Internet at
	http://www.illumos.org/license/CDDL.

	}}}

	Copyright 2020 OmniOS Community Edition (OmniOSce) Association.
-->

<service_bundle type="manifest" name="nagios">

<service
        name="ooce/application/nagios"
        type="service"
        version="1">

        <instance name='default' enabled='false'>

        <dependency
                name='network'
                grouping='require_all'
                restart_on="error"
                type="service">
                <service_fmri value='svc:/milestone/network:default' />
        </dependency>

        <dependency
                name="filesystem-local"
                grouping="require_all"
                restart_on="none"
                type="service">
                <service_fmri value="svc:/system/filesystem/local:default"/>
        </dependency>

        <dependency
                name='config-file'
                grouping='require_all'
                restart_on='refresh'
                type='path'>
                <service_fmri value='file://localhost/etc/opt/ooce/nagios/nagios.cfg' />
        </dependency>

        <exec_method
            type="method"
            name="start"
            exec="/lib/svc/method/application-nagios start"
            timeout_seconds="60">
            <method_context security_flags="aslr">
                <method_credential user="nagios" group="nagios" />
            </method_context>
        </exec_method>

        <exec_method
            type="method"
            name="stop"
            exec="/lib/svc/method/application-nagios stop"
            timeout_seconds="60"/>

        <property_group name="startd" type="framework">
                <propval name="duration" type="astring" value="contract"/>
        </property_group>

        </instance>

        <stability value="Unstable" />

        <template>
                <common_name>
                        <loctext xml:lang="C">
                        Nagios - Extremely powerful network monitoring system
                        </loctext>
                </common_name>
        </template>

</service>

</service_bundle>
