# build/fcgiwrap/files/fcgiwrap.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
	{{{ CDDL HEADER

	This file and its contents are supplied under the terms of the
	Common Development and Distribution License ("CDDL"), version 1.0.
	You may only use this file in accordance with the terms of version
	1.0 of the CDDL.

	A full copy of the text of the CDDL should have accompanied this
	source. A copy of the CDDL is also available via the Internet at
	http://www.illumos.org/license/CDDL.

	}}}

	Copyright 2020 OmniOS Community Edition (OmniOSce) Association.
-->

<service_bundle type='manifest' name='fcgiwrap'>

<service
	name='ooce/application/fcgiwrap'
	type='service'
	version='1'>

	<instance name='default' enabled='false'>

	<dependency
		name='loopback'
		grouping='require_any'
		restart_on='error'
		type='service'>
		<service_fmri value='svc:/network/loopback' />
	</dependency>

	<dependency
		name='network'
		grouping='optional_all'
		restart_on='error'
		type='service'>
		<service_fmri value='svc:/milestone/network' />
	</dependency>

	<dependency
	    name='filesystem_local'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/filesystem/local:default' />
	</dependency>

	<dependent
		name='fcgiwrap_multi-user'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri value='svc:/milestone/multi-user' />
	</dependent>

	<exec_method
	    type='method'
	    name='start'
	    exec='/opt/ooce/sbin/fcgiwrap -s %{config/socket} -c %{config/children}'
	    timeout_seconds='60'>
            <method_context security_flags="aslr">
		<method_credential user='fcgiwrap' group='fcgiwrap'
		    privileges='basic,!proc_info,!proc_session,!file_link_any' />
	    </method_context>
	</exec_method>

	<exec_method
	    type='method'
	    name='stop'
	    exec=':kill'
	    timeout_seconds='60' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='child' />
	</property_group>

	<property_group name='config' type='application'>
		<propval name='children' type='count' value='1' />
		<propval name='socket' type='astring'
		    value='unix:/var/opt/ooce/fcgiwrap/run/fcgiwrap.sock' />
	</property_group>

	</instance>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			fcgiwrap
			</loctext>
		</common_name>
	</template>

</service>

</service_bundle>

